# provision/etc/rc.local
#!/bin/bash



echo "Mounting $(vmware-hgfsclient) to /mnt"
vmhgfs-fuse .host:$(vmware-hgfsclient) /mnt -o allow_other -o noempty

echo "KickGuest VMware :: happy container
sleep 1

